# infra/igor/cef-runner.container
# ⚠️ Security warning
#
# This self-hosted runner is configured for a PUBLIC repository.
# Forks cannot trigger it because `build-vaapi.yml` rejects
# `pull_request` and `pull_request_target` events.
# Only push to main, scheduled, and manual dispatch trigger builds.
# Never add `pull_request` triggers to workflows using this runner.
Documentation=https://github.com/bnema/cef-vaapi
Wants=network-online.target
After=network-online.target podman.service

[Container]
Image=localhost/cef-runner:latest

# Environment variables (from .env file in same directory)
EnvironmentFile=%h/.config/containers/cef-runner.env

# Volumes
# Persistent work directory for runner state and build caches
Volume=%h/.local/share/cef-runner/_work:/home/runner/_work:Z
# Persistent makepkg/download cache (avoids re-downloading Chromium tarball)
Volume=cef-build-cache:/cache:Z

# Chromium build needs additional syscalls
SecurityOpt=seccomp=unconfined

# Resource limits
MemoryMax=28G
CPUQuota=1600%

# Auto-restart
AutoUpdate=local

[Service]
Restart=always
RestartSec=30

[Install]
WantedBy=default.target
